# gabinet/conf/production/nginx_supervisor.conf
[program:nginx]
directory=/usr/sbin/
autostart=true
autorestart=true
command=nginx
stdout_logfile=/var/www/logs/supervisor_nginx.log
redirect_stderr=true
